# dot-config/systemd##!docker/user/bluetooth-audio-default.service
[Unit]
Description=Make connected Bluetooth audio devices the default sink
After=pipewire.service pipewire-pulse.service wireplumber.service
Requires=pipewire-pulse.service

[Service]
Type=simple
ExecStart=%h/.local/bin/audio-control.sh watch-bluetooth
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
